# crontab
0 * * * * root /usr/local/bin/python /app/twitter_search.py >> /var/log/cron.log 2>&1
